# .config/systemd/user/hyprshell.service.d/override.conf
# Point hyprshell at the json5 config tracked in dotfiles (default path is
# config.ron, which hyprshell does NOT auto-substitute for config.json5).
[Service]
ExecStart=
ExecStart=/usr/bin/hyprshell run -c %h/.config/hyprshell/config.json5
